# roles/falcon_agent_install/templates/rabbitmq-monitor.cron.j2
# vim: set filetype=crontab:
*/1 * * * * root /opt/falcon-agent/rabbitmq-monitor.py 1>/dev/null
